# autostartx-rc
#!/sbin/runscript

depend() {
    need local
    use xfs
}

start() {
    ebegin "Starting X11 as ${X11USER}"
    start-stop-daemon --start --background --quiet --pidfile ${PIDFILE} --make-pidfile --chdir ${X11HOME} --exec /sbin/autostartx -- ${X11USER}
    eend $?
}

stop() {
    ebegin "Stopping X11 as ${X11USER}"
    start-stop-daemon --stop --quiet --oknodo --user ${X11USER} --pidfile ${PIDFILE} 
    eend $?
}
